# systemd/citadel-current-watcher.path
[Unit]
Description=Current realm directory watcher
Before=launch-default-realm.service

[Path]
PathChanged=/run/citadel/realms/current

[Install]
WantedBy=multi-user.target
