# systemd/citadel-desktop-watcher.service
[Unit]
Description=Desktop Watcher
StartLimitIntervalSec=0

[Service]
Type=oneshot
ExecStart=/usr/libexec/citadel-desktop-sync
